# files/etc/service/cron/run
#!/bin/sh

exec crond -f > /var/log/cron.log 2>&1
